# dynipv6.service
[Unit]
Description=dynipv6 - minimal self-hosted Dynamic DNS service
After=network.target

[Service]
Type=simple
# Runs as root because it writes the nginx config and reloads nginx.
# If you don't use the nginx feature, you can run it as an unprivileged user.
ExecStart=/opt/dynipv6/venv/bin/python3 /opt/dynipv6/dynipv6.py
Restart=always
RestartSec=10

# Hardening: only these paths are writable.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/dynipv6 /etc/nginx

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dynipv6

[Install]
WantedBy=multi-user.target
