# dynipv6-webui.service
[Unit]
Description=Dynamic IPv6/IPv4 DDNS Web Admin Panel
Documentation=file:///usr/share/doc/dynipv6/README.md
After=network.target dynipv6.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/dynipv6

# Start service
ExecStart=/opt/dynipv6/venv/bin/python3 /opt/dynipv6/webui.py

# Environment
Environment="FLASK_ENV=production"

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/dynipv6 /var/log/dynipv6 /etc/dynipv6

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dynipv6-webui

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
